# deploy/syncs3.timer
[Unit]
Description=Load synchronization with AWS S3 each 3 hours

[Timer]
OnBootSec=10s
OnUnitActiveSec=180m

[Install]
WantedBy=timers.target